# deploy/systemd/lumen-backup.service
[Unit]
Description=Lumen periodic backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
EnvironmentFile=-/opt/lumen/.env
ExecStart=/opt/lumen/scripts/backup.sh
StandardOutput=journal
StandardError=journal
