# deploy/systemd/lumen-storage-apply.service
[Unit]
Description=Lumen storage apply (admin-triggered remount cycle)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
Environment=LUMEN_STORAGE_TARGET=/opt/lumendata
ExecStart=/usr/local/sbin/lumen-storage-mount apply
# 1) Clean up the trigger so the next click re-fires PathChanged.
# 2) Restart every watcher whose inotify handle points inside the swapped
#    backup mount. A mount replacement can leave path units active but attached
#    to the old inode, so update/warm/backup/restore must all be re-armed.
ExecStopPost=/bin/sh -c 'rm -f /var/lib/lumen-storage/apply.trigger; for unit in lumen-update.path lumen-update-warm.path lumen-backup.path lumen-restore.path; do if ! systemctl try-restart "$$unit"; then echo "lumen-storage-apply: failed to restart $$unit" >&2; fi; done; exit 0'
TimeoutStartSec=600
