# deploy/systemd/lumen-tgbot.service
[Unit]
Description=Lumen Telegram Bot
After=network-online.target lumen-api.service
# Wants=network-online.target + lumen-api.service：软依赖，api 没起来 tgbot 启动会
# 失败但不会被 cascade 拉爆。避免 `Requires=lumen-api.service` 在 api crash-loop
# 时把 tgbot 一起拖进 StartLimitBurst 速死循环。
Wants=network-online.target lumen-api.service
PartOf=lumen-api.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=lumen
Group=lumen
WorkingDirectory=/opt/lumen/current/apps/tgbot
EnvironmentFile=-/opt/lumen/shared/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/lumen/current/.venv/bin/python -m app.main
# 保留 always：管理后台 restart 命令让进程 clean exit，需要 systemd 拉起并重载配置。
# disabled / 缺密钥时进程会保持 paused 且继续监听管理命令，不会触发重启或 start-limit。
Restart=always
RestartSec=3s
TimeoutStartSec=30s
TimeoutStopSec=20s
KillSignal=SIGTERM
KillMode=mixed
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
PrivateDevices=true
CapabilityBoundingSet=
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lumen-tgbot

[Install]
WantedBy=multi-user.target
