# deploy/systemd/lumen-update-runner.service
[Unit]
Description=Lumen one-click update runner
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/lumen
EnvironmentFile=-/opt/lumendata/backup/.update.env
Environment=LUMEN_UPDATE_NONINTERACTIVE=1
Environment=LUMEN_UPDATE_GIT_PULL=1
Environment=LUMEN_UPDATE_BUILD=1
StandardOutput=append:/opt/lumendata/backup/.update.log
StandardError=append:/opt/lumendata/backup/.update.log
ExecStart=/usr/bin/env bash /opt/lumen/scripts/update.sh
# update.sh 任何非零退出（git 冲突、依赖装失败、next build 失败等）都被当作
# success，避免 unit 进 failed 状态导致 path 单元下一次 PathChanged 不再触发
# StartUnit。真正的失败信息在 .update.log 里，前端 log_tail 直接能看到。
SuccessExitStatus=0 1 2 SIGTERM SIGINT
# Clean trigger / marker / env regardless of update result so the next click works.
ExecStopPost=/bin/sh -c "rm -f /opt/lumendata/backup/.update.trigger /opt/lumendata/backup/.update.running /opt/lumendata/backup/.update.env"
TimeoutStartSec=infinity
KillSignal=SIGTERM
KillMode=mixed
