# deploy/systemd/lumen-storage-test.service
[Unit]
Description=Lumen storage test (admin-triggered SMB connectivity test)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=LUMEN_STORAGE_TARGET=/opt/lumendata
ExecStart=/usr/local/sbin/lumen-storage-mount test
ExecStopPost=/bin/sh -c "rm -f /var/lib/lumen-storage/test.trigger"
TimeoutStartSec=120
