# deploy/systemd/lumen-storage-test.service
[Unit]
Description=Lumen storage test (admin-triggered SMB connectivity test)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/var/lib/lumen-storage/test.env
ExecStart=/usr/local/sbin/lumen-storage-mount test
ExecStopPost=/bin/sh -c "rm -f /var/lib/lumen-storage/test.trigger /var/lib/lumen-storage/test.env"
TimeoutStartSec=120
